# vault/server/vault-watcher.service
[Unit]
Description=Vault autocommit watcher (debounced git commit + push)
After=network-online.target vault-mcp.service
Wants=network-online.target

[Service]
Type=simple
User=misha
Group=misha
# WorkingDirectory is /opt/vault-mcp so `-m vault_watcher` finds vault_watcher.py.
WorkingDirectory=/opt/vault-mcp
Environment=VAULT_ROOT=/home/misha/vault
Environment=POLL_INTERVAL=5
Environment=DEBOUNCE_SECONDS=30
Environment=AUTO_PUSH=1
ExecStart=/opt/vault-mcp/venv/bin/python -m vault_watcher
Restart=on-failure
RestartSec=5s

# See note on vault-mcp.service: ProtectHome would block writes to /home/.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/misha/vault
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
