# vault/server/vault-mcp.service
[Unit]
Description=Vault MCP server (markdown vault over SSE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=misha
Group=misha
# WorkingDirectory is /opt/vault-mcp so `-m vault_mcp` finds vault_mcp.py.
# VAULT_ROOT (where the markdown vault lives) is set via env, NOT cwd.
WorkingDirectory=/opt/vault-mcp
Environment=VAULT_ROOT=/home/misha/vault
Environment=VAULT_BIND=0.0.0.0:8088
ExecStart=/opt/vault-mcp/venv/bin/python -m vault_mcp
Restart=on-failure
RestartSec=3s

# Lock it down a bit. The server only needs to touch the vault and write to git.
# NOTE: ProtectHome=read-only would block writes under /home/, which is exactly
# where the vault lives — leave it off and rely on ReadWritePaths + ProtectSystem.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/misha/vault
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
